# linux/systemd-cron/gentoo-package-check.timer
[Unit]
Description=Periodically check Portage tree age

[Timer]
OnActiveSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timer.target
