# linux/systemd-cron/gentoo-package-check.service
[Unit]
Description=Check Portage tree age

[Service]
Type=oneshot
ExecStart=/home/geoff/bin/gentoo-package-check

